# patches/haproxy.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<service_bundle type='manifest' name='haproxy'>
<service name='network/haproxy' type='service' version='1'>
<create_default_instance enabled='false' />
<single_instance />
<dependency name='network' grouping='require_all' restart_on='none' type='service'>
  <service_fmri value="svc:/milestone/network:default" />
</dependency>
<dependency name='fs-local' grouping='require_all' restart_on='none' type='service'>
  <service_fmri value='svc:/system/filesystem/local' />
</dependency>
<dependency name='config-file' grouping='require_all' restart_on='refresh' type='path'>
  <service_fmri value='file://localhost/etc/haproxy.cfg' />
</dependency>
<exec_method type='method' name='start'
	     exec='/opt/tribblix/haproxy/bin/haproxy-smf-wrapper /opt/tribblix/haproxy/bin/haproxy -p /var/haproxy/haproxy.pid -f /etc/haproxy.cfg' timeout_seconds='60'>
</exec_method>
<exec_method type='method' name='stop' exec=':kill' timeout_seconds='60' />
<exec_method type='method' name='refresh' exec=':kill -HUP' timeout_seconds='60' />
<property_group name='startd' type='framework'>
  <!-- sub-process core dumps shouldn't restart session -->
  <propval name='ignore_error' type='astring' value='core,signal' />
  <!-- the smf wrapper is a persistent child, not a daemon -->
  <propval name='duration' type='astring' value='child' />
</property_group>
<stability value='External' />
<template>
  <common_name>
    <loctext xml:lang='C'>haproxy server</loctext>
  </common_name>
</template>
</service>
</service_bundle>
