# patches/kibana.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<service_bundle type='manifest' name='kibana'>
<service name='network/kibana' type='service' version='1'>
<create_default_instance enabled='false' />
<single_instance />
<dependency name='fs-local' grouping='require_all' restart_on='none' type='service'>
  <service_fmri value='svc:/system/filesystem/local' />
</dependency>
<exec_method type='method' name='start' exec='/opt/tribblix/kibana/bin/kibana' timeout_seconds='60'>
  <method_context working_directory='/opt/tribblix/kibana'>
    <method_credential user='kibana' group='kibana' />
  </method_context>
</exec_method>
<exec_method type='method' name='stop' exec=':kill' timeout_seconds='60' />
<property_group name='startd' type='framework'>
  <!-- sub-process core dumps shouldn't restart session -->
  <propval name='ignore_error' type='astring' value='core,signal' />
  <!-- not a daemonized service -->
  <propval name='duration' type='astring' value='child' />
</property_group>
<stability value='External' />
<template>
  <common_name>
    <loctext xml:lang='C'>Kibana server</loctext>
  </common_name>
</template>
</service>
</service_bundle>
